# supervisord.conf
[supervisord]
nodaemon=true

[program:swoole]
command=php bin/console swoole:http start
autostart=true
autorestart=true
stderr_logfile=/var/log/swoole.err.log
stdout_logfile=/var/log/swoole.out.log